# scripts/httpserver.conf
[program:http_ci]
command=python /opt/cmc.website/deploy/scripts/http_supervisor.py
directory=/opt/cmc.website/deploy/scripts
numprocs=1
numprocs_start=0
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s-%(process_num)01d.log
user=root
